# pi-conf/telegrambot.service
[Unit]
Description=Heating Telegram Bot

[Service]
Type=simple
ExecStart=/bin/bash /home/pi/bin/telegrambot.sh

[Install]
WantedBy=multi-user.target
